# systemd/kyaclaw-monitor.service
[Unit]
Description=KyaClaw PRISM File Integrity Monitor
After=network.target

[Service]
Type=simple
User=__KYACLAW_USER__
ExecStart=/usr/bin/node --enable-source-maps /opt/openclaw-prism/packages/monitor/dist/index.js
Environment=KYACLAW_MONITOR_START=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
