# systemd/kyaclaw-proxy.service
[Unit]
Description=KyaClaw PRISM Invoke Guard Proxy
After=network.target kyaclaw-scanner.service
Wants=kyaclaw-scanner.service

[Service]
Type=simple
User=__KYACLAW_USER__
ExecStart=/usr/bin/node --enable-source-maps /opt/openclaw-prism/packages/proxy/dist/index.js
Environment=INVOKE_GUARD_POLICY=/opt/openclaw-prism/config/invoke-guard.policy.json
Environment=KYACLAW_PROXY_START=1
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
