# launchd/com.kyaclaw.scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kyaclaw.scanner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--enable-source-maps</string>
    <string>/opt/openclaw-prism/packages/scanner/dist/index.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SCANNER_HOST</key>
    <string>127.0.0.1</string>
    <key>SCANNER_PORT</key>
    <string>18766</string>
    <key>KYACLAW_SCANNER_START</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/kyaclaw-scanner.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/kyaclaw-scanner.err</string>
</dict>
</plist>
