# launchd/com.kyaclaw.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kyaclaw.monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--enable-source-maps</string>
    <string>/opt/openclaw-prism/packages/monitor/dist/index.js</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENCLAW_AUDIT_HMAC_KEY</key>
    <string>replace-with-secure-random-key</string>
    <key>KYACLAW_MONITOR_START</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/kyaclaw-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/kyaclaw-monitor.err</string>
</dict>
</plist>
